# volumes/supervisord.conf
[supervisord]
nodaemon=true
[program:nrpe]
command=/bin/sh -c "nrpe -c /etc/nagios/nrpe.cfg -d"
startsecs=0
exitcodes=1
autorestart=false
